# deploy/systemd/nuwax-host-gateway.service.example
[Unit]
Description=nuwax_robot_bridge 宿主机网关
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=unitree
Group=unitree
WorkingDirectory=/home/unitree/nuwax-robot-bridge
# 默认会读取 WorkingDirectory 下的 .env。
# 如果你想把配置独立放到 /etc，再取消下一行注释并准备对应文件。
# EnvironmentFile=/etc/nuwax_robot_bridge/gateway.env
# 平台文件日志路径由 .env 中的 NUWAX_LOG_* 控制。
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 -m uvicorn gateways.app:create_gateway_app_from_env --factory --host 0.0.0.0 --port 8766
Restart=always
RestartSec=3
TimeoutStopSec=20

# 建议通过宿主机目录持久化日志和制品
ReadWritePaths=/var/log/nuwax_robot_bridge /var/lib/nuwax_robot_bridge

[Install]
WantedBy=multi-user.target
