# deploy/systemd/nuwax-relay.service.example
[Unit]
Description=nuwax_robot_bridge 容器侧 relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=nuwax
Group=nuwax
WorkingDirectory=/home/ml/Deeprob/nuwax_robot_bridge
EnvironmentFile=/etc/nuwax_robot_bridge/relay.env
ExecStart=/usr/bin/python3 -m uvicorn gateways.app:create_relay_app_from_env --factory --host 0.0.0.0 --port 9766
Restart=always
RestartSec=3
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
